# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:mongod]
command=/usr/bin/mongod

[program:mongo_stup]
command=./setup_mongo_user.sh
